# dev-watch.service
[Unit]
Description=dev-watch — local process monitor
After=network.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__PATH__
ExecStart=__PATH__/.venv/bin/python3 -m src.server
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
